# modules/systemd/templates/mnt-target.mount.erb
[Unit]
Description=<%= @name %>

[Mount]
What=<%= @what %>
Where=/mnt/<%= @name %>
Type=<%= @type %>

[Install]
WantedBy=multi-user.target
