# modules/systemd/templates/mnt-target.automount.erb
[Unit]
Description=Automount <%= @name %>

[Automount]
Where=/mnt/<%= @name %>

[Install]
WantedBy=multi-user.target
